# scripts/ffmpeg-capture.service
[Unit]
Description=Soundcard Capture -> Icecast (lossless FLAC)
Documentation=https://github.com/FreieRadios/radiobox
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=tsing
Group=audio
EnvironmentFile=/path/to/.env
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/ffmpeg -hide_banner -loglevel info \
   -f alsa -ac 2 -ar 48000 -i hw:CARD=UGM96,DEV=0 \
   -c:a flac -compression_level 5 \
   -content_type application/ogg \
   -f ogg \
   icecast://source:${ICECAST_SOURCE_PASSWORD}@${ICECAST_HOSTNAME}:${ICECAST_PORT}/lossless-capture
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ffmpeg-capture

[Install]
WantedBy=multi-user.target
